# scripts/common/profile.d/99nvidia-wsl.sh
#!/usr/bin/env sh

# if [ -n "$WSLENV" ]; then
#   export PATH="/usr/lib/wsl/lib:${PATH}"
# fi
