# scripts/common/profile.d/less-pager.sh
#!/usr/bin/env sh

export LESS='-Rf'
export PAGER=/usr/bin/less
