# scripts/common/profile.d/tmux-conf.sh
#!/usr/bin/env sh

export ZSH_TMUX_CONFIG="${XDG_CONFIG_HOME:-${HOME}/.config}/tmux/tmux.conf"
export TMUX_PLUGIN_MANAGER_PATH="${XDG_CONFIG_HOME:-${HOME}/.config}/tmux/plugins"
